# files/root/etc/supervisor/conf.d/mysql.conf
[program:mysql]
command=/usr/bin/mysqld_safe --initialize-insecure 
autorestart=true
